# telegraf/zpool_influxdb.timer
[Unit]
Description=Runs zpool_influxdb

[Timer]
Unit=zpool_influxdb.service
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
